# app-misc/keyd-init/files/keyd.init
#!/sbin/openrc-run

name="$RC_SVCNAME"
description="key remapping daemon"
command="/usr/bin/keyd"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	after localmount
	need localmount
}
