# sys-apps/files/kmscon.init
#!/sbin/openrc-run

name="$RC_SVCNAME"
description="Linux KMS/DRM based virtual Console Emulator"
command="/usr/bin/kmscon"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
}
